# debian/postinst
#!/bin/sh

set -e

case "$1" in
    configure)
        invoke-rc.d acpid restart>/dev/null
        if [ -x /etc/alsa/modprobe-post-install.d/eeepc-acpi-scripts ]; then
          /etc/alsa/modprobe-post-install.d/eeepc-acpi-scripts
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "preinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

#DEBHELPER#

exit 0

